# systemd/skipa-banlist.service
# Описание unit-файла systemd.
[Unit]
# Человекочитаемое описание сервиса.
Description=Refresh UFW SKIPA blocklist from CyberOK_Skipa_ips
# Запускать сервис только после того, как сеть считается поднятой.
After=network-online.target
# По возможности подтягивать network-online.target вместе с этим сервисом.
Wants=network-online.target

# Описание параметров самого сервиса.
[Service]
# Сервис однократный: он выполняет updater-скрипт и завершается.
Type=oneshot
# Основная команда сервиса: скачать свежий список и обновить managed-блок в /etc/ufw/before.rules.
ExecStart=/usr/local/sbin/update-skipa-banlist.sh

# Описание секции установки unit-файла.
[Install]
# Сервис обычно стартует через timer, но его можно запускать и вручную.
WantedBy=multi-user.target
